# contrib/openrc/zenbook-kb-hotkeys
#!/sbin/openrc-run
# OpenRC service: listen for KEY_KBDILLUM* and call kb-brightness.
#
# Install via configure.py (writes /etc/conf.d/zenbook-kb-hotkeys with command_user):
#   sudo python3 configure.py --defaults --all-yes
#
# Do not copy this file by hand — install via configure.py (writes /etc/conf.d/).

name="zenbook-kb-hotkeys"
description="Zenbook Duo Fn+ / special-key handler"

: "${command:=/usr/local/bin/kb-brightness-hotkeys}"
: "${command_background:=yes}"
: "${command_args:=}"
: "${pidfile:=/run/${RC_SVCNAME}.pid}"
: "${output_log:=/var/log/${RC_SVCNAME}.log}"
: "${error_log:=/var/log/${RC_SVCNAME}.log}"

depend() {
	need localmount
	after hotplug
}

start_pre() {
	# command_user is set in /etc/conf.d/zenbook-kb-hotkeys by configure.py.
	if [ -z "${command_user:-}" ]; then
		eerror "command_user is not set — run: sudo python3 configure.py --defaults --all-yes"
		return 1
	fi
	# OpenRC only appends to an existing log; create it before the first start.
	if [ -n "${output_log:-}" ]; then
		checkpath -f -m 0664 -o "$command_user" "$output_log"
	fi
}
