# contrib/openrc/zenbook-screenpad-sync
#!/sbin/openrc-run
# Mirror main panel brightness onto ASUS ScreenPad Plus.

name="zenbook-screenpad-sync"
description="Zenbook ScreenPad brightness sync"

: "${command:=/usr/bin/screenpad-sync}"
: "${command_args:=--replace}"
: "${command_background:=yes}"
: "${pidfile:=/run/${RC_SVCNAME}.pid}"
: "${output_log:=/var/log/${RC_SVCNAME}.log}"
: "${error_log:=/var/log/${RC_SVCNAME}.log}"

depend() {
	need localmount
	after zenbook-screenpad modules
}

start_pre() {
	for _wait_lib in \
		/usr/share/zenbook-scripts/lib/openrc-wait.sh \
		/usr/local/share/zenbook-scripts/lib/openrc-wait.sh; do
		if [ -f "${_wait_lib}" ]; then
			# shellcheck source=lib/openrc-wait.sh
			. "${_wait_lib}"
			zenbook_rc_clear_stale_svc_pidfile
			break
		fi
	done
	if [ -x /usr/local/bin/screenpad-boot ]; then
		SCREENPAD_BOOT_MODE=on /usr/local/bin/screenpad-boot >/dev/null 2>&1 || true
	elif [ -x /usr/bin/screenpad-boot ]; then
		SCREENPAD_BOOT_MODE=on /usr/bin/screenpad-boot >/dev/null 2>&1 || true
	fi
	if [ -n "${output_log:-}" ]; then
		checkpath -f -m 0644 -o root:root "${output_log}"
	fi
}
