# contrib/openrc/zenbook-platform-fan-control
#!/sbin/openrc-run
# Adaptive platform fan-profile daemon (platform-fan-control).
#
# Machine-global config (not per-user):
#   cp /usr/share/zenbook-scripts/fan-control.json.example \
#      /etc/zenbook-scripts/fan-control.json
#   sudo rc-update add zenbook-platform-fan-control default
#   sudo rc-service zenbook-platform-fan-control start
#
# Override path with FAN_CONTROL_CONFIG in /etc/conf.d/zenbook-platform-fan-control.

name="zenbook-platform-fan-control"
description="Adaptive laptop fan / platform_profile controller"

: "${command:=/usr/bin/platform-fan-control}"
: "${command_args:=run}"
: "${command_background:=yes}"
: "${pidfile:=/run/${RC_SVCNAME}.pid}"
: "${output_log:=/var/log/${RC_SVCNAME}.log}"
: "${error_log:=/var/log/${RC_SVCNAME}.log}"
: "${FAN_CONTROL_CONFIG:=/etc/zenbook-scripts/fan-control.json}"

depend() {
	need localmount
	use logger
	after modules
}

start_pre() {
	if [ ! -x "${command}" ]; then
		for c in /usr/bin/platform-fan-control /usr/local/bin/platform-fan-control \
			/usr/bin/kb-fan-control /usr/local/bin/kb-fan-control; do
			if [ -x "${c}" ]; then
				command="${c}"
				break
			fi
		done
	fi
	if [ ! -x "${command}" ]; then
		eerror "platform-fan-control not found"
		return 1
	fi
	cfg="${FAN_CONTROL_CONFIG}"
	if [ ! -f "${cfg}" ]; then
		ewarn "missing ${cfg}"
		ewarn "copy fan-control.json.example → /etc/zenbook-scripts/fan-control.json"
	fi
	command_args="run -c ${cfg}"
	if [ -n "${output_log:-}" ]; then
		checkpath -f -m 0644 -o root:root "${output_log}"
	fi
}
