# contrib/openrc/zenbook-kb-hotkeys
#!/sbin/openrc-run
# OpenRC service: listen for KEY_KBDILLUM* and call kb-brightness.
#
# Install via configure.py (writes /etc/conf.d/zenbook-kb-hotkeys with command_user):
#   sudo python3 configure.py --defaults --all-yes
#
# Do not copy this file by hand — install via configure.py (writes /etc/conf.d/).

name="zenbook-kb-hotkeys"
description="Zenbook Duo Fn+ / special-key handler"

: "${command:=/usr/bin/kb-brightness-hotkeys}"
: "${command_background:=yes}"
: "${command_args:=--replace}"
: "${pidfile:=/run/${RC_SVCNAME}.pid}"
: "${output_log:=/var/log/${RC_SVCNAME}.log}"
: "${error_log:=/var/log/${RC_SVCNAME}.log}"

depend() {
	need localmount
	after hotplug zenbook-kb-hid-asus
}

start_pre() {
	for _wait_lib in \
		/usr/share/zenbook-scripts/lib/openrc-wait.sh \
		/usr/local/share/zenbook-scripts/lib/openrc-wait.sh; do
		if [ -f "${_wait_lib}" ]; then
			# shellcheck source=lib/openrc-wait.sh
			. "${_wait_lib}"
			zenbook_rc_clear_stale_svc_pidfile
			zenbook_rc_wait_not_starting zenbook-kb-hid-asus 20 || true
			break
		fi
	done
	# command_user is set in /etc/conf.d/zenbook-kb-hotkeys by configure.py.
	if [ -z "${command_user:-}" ]; then
		eerror "command_user is not set — run: sudo python3 configure.py --defaults --all-yes"
		return 1
	fi
	# Avoid flock races when zenbook-kb-hid-asus stops/restarts us during sideload.
	# Do not poll rc-service on *this* service — start_pre already runs in "starting".
	if [ -f /run/zenbook-kb-hotkeys.pid ]; then
		i=0
		while kill -0 "$(cat /run/zenbook-kb-hotkeys.pid 2>/dev/null)" 2>/dev/null; do
			sleep 0.2
			i=$((i + 1))
			[ "$i" -ge 15 ] && break
		done
	fi
	# Drop orphaned listener left after a partial start (avoids "already running").
	if command -v start-stop-daemon >/dev/null 2>&1; then
		start-stop-daemon --stop --exec "${command}" --retry TERM/3/KILL/2 \
			>/dev/null 2>&1 || true
		rm -f /run/zenbook-kb-hotkeys.pid
	fi
	# OpenRC only appends to an existing log; create it before the first start.
	if [ -n "${output_log:-}" ]; then
		checkpath -f -m 0664 -o "$command_user" "$output_log"
	fi
}

start_post() {
	for _sleep in /usr/bin/kb-brightness-sleep /usr/local/bin/kb-brightness-sleep; do
		if [ -x "${_sleep}" ]; then
			"${_sleep}" restore 2>/dev/null || true
			break
		fi
	done
}

stop_pre() {
	for _sleep in /usr/bin/kb-brightness-sleep /usr/local/bin/kb-brightness-sleep; do
		if [ -x "${_sleep}" ]; then
			"${_sleep}" save 2>/dev/null || true
			break
		fi
	done
}

stop() {
	# Quiet no-op when already stopped (avoids start-stop-daemon spam).
	_pid=""
	if [ -n "${pidfile:-}" ] && [ -f "${pidfile}" ]; then
		_pid="$(tr -d '[:space:]' < "${pidfile}" 2>/dev/null || true)"
	fi
	if [ -n "${_pid}" ] && kill -0 "${_pid}" 2>/dev/null; then
		default_stop
		return $?
	fi
	rm -f "${pidfile:-}"
	# Orphan without pidfile (start_pre may have left one).
	if command -v pgrep >/dev/null 2>&1 && [ -n "${command:-}" ]; then
		if pgrep -x "$(basename "${command}")" >/dev/null 2>&1; then
			default_stop
			return $?
		fi
	fi
	ebegin "Stopping ${RC_SVCNAME:-zenbook-kb-hotkeys}"
	eend 0
	return 0
}
