# contrib/udev/99-zenbook-screenpad.rules
# ASUS ScreenPad Plus backlight — allow video group to control brightness/power
ACTION=="add|change", SUBSYSTEM=="backlight", KERNEL=="asus_screenpad", \
  RUN+="/bin/chgrp video /sys/class/backlight/%k/brightness /sys/class/backlight/%k/bl_power", \
  RUN+="/bin/chmod g+w /sys/class/backlight/%k/brightness /sys/class/backlight/%k/bl_power"
